# etc/pf/pf.openwrt
#!/bin/sh /etc/rc.common

START=19
USE_PROCD= # 0

restart() {
  start
}

start() {
  cd /etc/pf
  exec ./pf start
}

stop() {
  cd /etc/pf
  exec ./pf stop
}

EXTRA_HELP="	load	alias for start
	status	shows current rules
	stat	alias for status
"

EXTRA_COMMANDS="load status stat"

load() { start; }
status() { iptables-save; }
stat() { status; }

