# etc/ram-pkg.list
# config file for /etc/init.d/ram-pkg
# for local packages installed in RAM available after uplink is avail
#
#- debugging aid
#strace
#lsof
#diffutils

#- collectd
#collectd
#collectd-mod-conntrack
#collectd-mod-df
#collectd-mod-disk
#collectd-mod-entropy
#collectd-mod-exec
#collectd-mod-interface
#collectd-mod-iwinfo
#collectd-mod-load
#collectd-mod-memory
#collectd-mod-ping
#collectd-mod-processes
#collectd-mod-syslog
#collectd-mod-thermal
#collectd-mod-uptime
#collectd-mod-vmem
#collectd-mod-wireless
#collectd-mod-write-graphite

# haveged - randomness daemon
#haveged

#- pkill/pgrep - NEED: ncurses
#procps-ng-pkill
#procps-ng-ps

#- htop - NEED: ncurses
#htop

#- vim - NEED: ncurses
#vim
#vim-full
