# etc/ram-pkg.list
# config file for /etc/init.d/ram-pkg
# for local packages installed in RAM available after uplink is avail
#
#- debugging aid
#strace
#lsof
#diffutils

#- pkill/pgrep
#procps-ng-pkill
#procps-ng-ps

#- htop
#htop

#- vim
#vim
#vim-full

#- haveged - randomness daemon
#haveged

#- basic LuCi
#luci

#- collectd
#collectd-mod-syslog
#collectd-mod-write-graphite
#collectd-mod-df
#collectd-mod-disk
#collectd-mod-entropy
#collectd-mod-exec
#collectd-mod-interface
#collectd-mod-load
#collectd-mod-memory
#collectd-mod-processes
#collectd-mod-cpu
#collectd-mod-iwinfo
#collectd-mod-conntrack
#collectd-mod-ping
#collectd-mod-contextswitch
#collectd-mod-vmem
#collectd-mod-table
#collectd-mod-tail
#collectd-mod-uptime
