# etc/rc.local
#!/bin/sh
#  /etc/rc.local

#- user config files for root in temp legen
  cp -a /root/.[a-z]* /tmp

#- resize tmp
  /usr/local/sbin/tmp-resize 10240

# end
exit 0
